# omsBackend/daphnew.service
[Unit]
Description=Daphne Interface new Server

[Service]
ExecStart=/bin/bash -c 'cd /data/projects/oms/omsBackend && /root/.pyenv/shims/python manage.py runworker --threads 5 --only-channels=websocket.*'
NonBlocking=true

[Install]
WantedBy=multi-user.target