# omsBackend/daphnew.service
[Unit]
Description=Daphne Interface Server
[Service]
ExecStartPre=/usr/bin/bash -c 'cd /data/projects/oms/omsBackend'
ExecStart=/usr/bin/bash -c '/root/.pyenv/shims/python manage.py runworker --threads 5 --only-channels=websocket.*'
Restart=always
KillSignal=SIGQUIT
Type=notify
NotifyAccess=all

[Install]
WantedBy=multi-user.target